# crontab
*/5 * * * * /usr/local/bin/node /app/dist/index.js >> /var/log/cron.log 2>&1
