# crontab
*/5 * * * * /usr/local/bin/node /app/dist/index.js > /dev/stdout
